# packaging/examples/systemd/ntip-online-backup.timer
[Unit]
Description=Example schedule for operator-managed NTIP SQLite backups

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true
Unit=ntip-online-backup.service

[Install]
WantedBy=timers.target
